# ansible/roles/minecraft/templates/root.crontab.j2
0 0,4,8,12,16,20 * * * /usr/sbin/service minecraft restart
